# assets/rootfs/usr/lib/systemd/system/systemd-journald.service.d/nice.conf
# Set Nice=-1 to dodge watchdog on soft lockups (LP: #1696970).

[Service]
Nice=-1
