# assets/rootfs/var/lib/dpkg/info/systemd-cryptsetup.list
/.
/usr
/usr/bin
/usr/bin/systemd-cryptenroll
/usr/bin/systemd-cryptsetup
/usr/lib
/usr/lib/nvpcr
/usr/lib/nvpcr/cryptsetup.nvpcr
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/cryptsetup-pre.target
/usr/lib/systemd/system/cryptsetup.target
/usr/lib/systemd/system/initrd-root-device.target.wants
/usr/lib/systemd/system/integritysetup-pre.target
/usr/lib/systemd/system/integritysetup.target
/usr/lib/systemd/system/remote-cryptsetup.target
/usr/lib/systemd/system/remote-integritysetup.target
/usr/lib/systemd/system/remote-veritysetup.target
/usr/lib/systemd/system/sysinit.target.wants
/usr/lib/systemd/system/system-systemd\x2dcryptsetup.slice
/usr/lib/systemd/system/system-systemd\x2dveritysetup.slice
/usr/lib/systemd/system/veritysetup-pre.target
/usr/lib/systemd/system/veritysetup.target
/usr/lib/systemd/system-generators
/usr/lib/systemd/system-generators/systemd-cryptsetup-generator
/usr/lib/systemd/system-generators/systemd-integritysetup-generator
/usr/lib/systemd/system-generators/systemd-veritysetup-generator
/usr/lib/systemd/systemd-integritysetup
/usr/lib/systemd/systemd-veritysetup
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/cryptsetup
/usr/lib/x86_64-linux-gnu/cryptsetup/libcryptsetup-token-systemd-fido2.so
/usr/lib/x86_64-linux-gnu/cryptsetup/libcryptsetup-token-systemd-pkcs11.so
/usr/lib/x86_64-linux-gnu/cryptsetup/libcryptsetup-token-systemd-tpm2.so
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/systemd-cryptenroll
/usr/share/doc
/usr/share/doc/systemd-cryptsetup
/usr/share/doc/systemd-cryptsetup/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/systemd-cryptsetup
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/systemd-cryptenroll.1.gz
/usr/share/man/man5
/usr/share/man/man5/integritytab.5.gz
/usr/share/man/man5/veritytab.5.gz
/usr/share/man/man8
/usr/share/man/man8/systemd-cryptsetup-generator.8.gz
/usr/share/man/man8/systemd-cryptsetup.8.gz
/usr/share/man/man8/systemd-integritysetup-generator.8.gz
/usr/share/man/man8/systemd-integritysetup@.service.8.gz
/usr/share/man/man8/systemd-veritysetup-generator.8.gz
/usr/share/man/man8/systemd-veritysetup@.service.8.gz
/usr/lib/systemd/system/initrd-root-device.target.wants/remote-cryptsetup.target
/usr/lib/systemd/system/initrd-root-device.target.wants/remote-integritysetup.target
/usr/lib/systemd/system/initrd-root-device.target.wants/remote-veritysetup.target
/usr/lib/systemd/system/sysinit.target.wants/cryptsetup.target
/usr/lib/systemd/system/sysinit.target.wants/integritysetup.target
/usr/lib/systemd/system/sysinit.target.wants/veritysetup.target
/usr/lib/systemd/systemd-cryptsetup
/usr/share/doc/systemd-cryptsetup/NEWS.Debian.gz
/usr/share/doc/systemd-cryptsetup/changelog.Debian.gz
/usr/share/man/man8/systemd-cryptsetup@.service.8.gz
/usr/share/man/man8/systemd-integritysetup.8.gz
/usr/share/man/man8/systemd-veritysetup.8.gz
